# qqbot-codex-bridge.service
[Unit]
Description=QQBot Codex Bridge (Direct ACP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/qqbot-codex-bridge
EnvironmentFile=/home/pi/qqbot-codex-bridge/.env
ExecStart=/usr/bin/node /home/pi/qqbot-codex-bridge/index.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
